# ex27_daemon/ex27_daemon.service
[Unit]
Description=Otus daemon.

[Service]
Type=simple
ExecStart=/usr/bin/ots_daemon

[Install]
WantedBy=multi-user.target
